# systemd/jam-ws-commands-announced.path
[Unit]
Description=Watch for device announcement to trigger WebSocket commands service

[Path]
# Watch for the announced flag (created by jam-announce)
PathExists=/etc/jam/device_data/.announced

# Trigger the service when the file appears
Unit=jam-ws-commands-trigger.service

[Install]
WantedBy=multi-user.target
