# systemd/jam-player-display.service
[Unit]
Description=JAM Player Display - Unified display manager for all display states
Documentation=https://github.com/anthropics/jam-player

# Wait for display manager (lightdm) to be ready before starting.
# Using display-manager.service instead of graphical.target avoids ordering cycles
# while still ensuring X11 is available.
After=display-manager.service

# Restart limits
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=notify

# Run as root for display access
User=root
Group=root

# Working directory
WorkingDirectory=/opt/jam/services

# Environment - DISPLAY and XAUTHORITY are required for MPV/feh to render
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/comitup/.Xauthority
Environment=PYTHONUNBUFFERED=1

# Wait for X display to be ready (fixes race condition where service starts before Xwayland)
# Simple sleep to allow lightdm/Xwayland to fully initialize
ExecStartPre=/bin/sleep 3

# Run the unified display manager (handles all 4 display modes)
ExecStart=/opt/jam/venv/bin/python /opt/jam/services/jam_player_display.py

# Clean up any display processes on stop
ExecStopPost=-/usr/bin/pkill -f 'feh.*jam_display'
ExecStopPost=-/usr/bin/pkill -f 'mpv.*mpv-socket'

# Always restart - display must always be running
Restart=always
RestartSec=5

# Watchdog - service must ping within 60 seconds
WatchdogSec=60

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jam-player-display

[Install]
WantedBy=multi-user.target
