# systemd/jam-ws-commands.service
[Unit]
Description=JAM Player WebSocket Commands - Receives real-time commands from backend
Documentation=https://github.com/anthropics/jam-player

# Only run when device has announced itself to the backend (has a DB record)
ConditionPathExists=/etc/jam/device_data/.announced

# Requires first boot to have completed (credentials must exist)
Requires=jam-first-boot.service
After=jam-first-boot.service

# Requires network and synced clock (for SSL certificate validation)
After=network-online.target time-sync.target
Wants=network-online.target time-sync.target

# Restart limits
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=notify

# Run as root (needed for network configuration changes)
User=root
Group=root

# Working directory
WorkingDirectory=/opt/jam/services

# Main script - maintains WebSocket connection and handles commands
ExecStart=/opt/jam/venv/bin/python /opt/jam/services/jam_ws_commands.py

# Always restart - commands are important for remote management
Restart=always
RestartSec=30

# Watchdog - must ping within 120 seconds (we ping every 30s)
WatchdogSec=120

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=JAM_ENVIRONMENT=production

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jam-ws-commands

[Install]
WantedBy=multi-user.target
